# supervisor/conf.d/radiusd.conf
[program:radiusd]
command=/usr/sbin/radiusd -f
autorestart=true
