# docker/supervisor/conf.d/pptpd.conf
[program:pptpd]
command=/usr/bin/pidproxy /var/run/pptpd.pid /usr/sbin/pptpd -f -c /etc/pptpd.conf
stdout_logfile=/var/log/pptpd.log
stderr_logfile=/var/log/pptpd.log
autorestart=true
